# web/rusty-proxy-again/src/supervisord.conf
[supervisord]
nodaemon=true

[program:proxy]
command=/home/chall/http_proxy/target/debug/http_proxy

[program:api]
directory=/home/chall/app/api
command=npm start